# onnxasr/rootfs/etc/services.d/onnxasr/run
#!/bin/bash

CACHE_DIR="/share/onnxasr"
export HOME="${CACHE_DIR}"
mkdir -p "${CACHE_DIR}"

export PYTHONPATH=/app:$PYTHONPATH

echo "[INFO] Starting ONNX ASR server..."

MODEL=$(jq -r '.model' /data/options.json)
INT8_QUANTIZATION=$(jq -r '.int8_quantization' /data/options.json)
DEVICE=$(jq -r '.device' /data/options.json)
PORT=$(jq -r '.port' /data/options.json)
DEBUG=$(jq -r '.debug' /data/options.json)

ARGS=(
    --model "${MODEL}"
    --uri "tcp://0.0.0.0:${PORT}"
    --device "${DEVICE}"
)

if [ "${INT8_QUANTIZATION}" = "true" ]; then
    ARGS+=(--quantization int8)
    echo "[INFO] Using int8 quantization."
else
    echo "[INFO] Running without specific quantization (using model default)."
fi

if [ "${DEBUG}" = "true" ]; then
    ARGS+=(--debug)
fi

echo "[INFO] Running with model: ${MODEL}, port: ${PORT}, device: ${DEVICE}, debug: ${DEBUG}"

exec python3 -m wyoming_onnxasr "${ARGS[@]}"
