# onnxasr/rootfs/etc/services.d/onnxasr/run
#!/bin/bash

# Определяем, где будут храниться модели
CACHE_DIR="/share/onnxasr"
export HOME="${CACHE_DIR}"
mkdir -p "${CACHE_DIR}"

# Добавляем /app в PYTHONPATH
export PYTHONPATH=/app:$PYTHONPATH

echo "[INFO] Starting ONNX ASR server..."

# Читаем опции из конфигурации аддона
MODEL=$(jq -r '.model // "vosk-model-small-ru"' /data/options.json)
QUANTIZATION=$(jq -r '.quantization // "int8"' /data/options.json)
DEVICE=$(jq -r '.device // "cpu"' /data/options.json)
PORT=$(jq -r '.port // 10305' /data/options.json)
DEBUG=$(jq -r '.debug // false' /data/options.json)

# Собираем аргументы для команды
ARGS=(
    --model "${MODEL}"
    --uri "tcp://0.0.0.0:${PORT}"
    --quantization "${QUANTIZATION}"
    --device "${DEVICE}"
)

# Добавляем флаг --debug, если включен
if [ "${DEBUG}" = "true" ]; then
    ARGS+=(--debug)
fi

echo "[INFO] Running ONNX ASR with model: ${MODEL}, port: ${PORT}, quantization: ${QUANTIZATION}, device: ${DEVICE}, debug: ${DEBUG}"

# Запускаем сервер
exec python3 -m wyoming_onnxasr "${ARGS[@]}" || {
    echo "[ERROR] Failed to start ONNX ASR server"
    exit 1
}
